# templ/CT0/etc/cron.d/syncbackup
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

17   0   *  *  * root   /etc/scripts/sync2backupsrv.sh > /dev/null 2>&1

